# BCCVL_Visualiser/buildout.d/supervisor.conf.in
[program:visualiser]
command = ${buildout:bin-directory}/pserve ${conf:file}
process_name = visualiser
directory = ${buildout:directory}
priority = 20
redirect_stderr = true
user = visualiser
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_logfile = ${buildout:directory}/var/log/visualiser
